# script/deploy/nubase.service
[Unit]
Description=Nubase backend (API + bundled Studio UI)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/nubase
ExecStart=/bin/bash /root/nubase/script/nubase_start.sh
Restart=always
RestartSec=10
# java exits with 143 on SIGTERM — treat as a clean stop, not a crash.
SuccessExitStatus=143
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
